# infra/systemd/tecnoria-chat-portal.service
[Unit]
Description=Tecnoria Chat Portal
After=network.target tecnoria-chat-api.service tecnoria-chat-widget.service
Wants=tecnoria-chat-api.service tecnoria-chat-widget.service

[Service]
Type=simple
User=tecnoria
Group=tecnoria
WorkingDirectory=/var/www/talkaris/apps/portal
EnvironmentFile=/etc/tecnoria/portal.env
ExecStart=/var/www/bin/run-release-entry.sh /var/www/talkaris/apps/portal/releases/current/server/server.mjs
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
